# launchd/com.codex-computer-use-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex-computer-use-bridge</string>

  <!-- Edit the path below to where you cloned this repo. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>/REPLACE/WITH/ABS/PATH/codex-computer-use-bridge/src/server.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/REPLACE/WITH/ABS/PATH/codex-computer-use-bridge</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/codex-computer-use-bridge.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/codex-computer-use-bridge.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>BRIDGE_HOST</key>
    <string>127.0.0.1</string>
    <key>BRIDGE_PORT</key>
    <string>37321</string>
  </dict>
</dict>
</plist>
